# docker/cron/jobs.crontab
RUST_LOG=debug
*/5 * * * * /app/release/onsen-compo comments process-image >/proc/1/fd/1 2>/proc/1/fd/2
